# ansible/roles/spark/templates/spark-history.service.j2
[Unit]
Description=Apache Spark Master 
After=network.target

[Service]
User={{ spark_user }}
Group={{ application_group }}

Environment=JAVA_HOME={{ java_home_17 }}
Environment=HADOOP_HOME={{ hadoop_home }}
Environment=HADOOP_CLASSPATH=$($HADOOP_HOME/bin/hadoop classpath)

ExecStart={{ spark_home }}/run-spark-history.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
